# package/openwrt/hostile/files/hostile.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 HostileWrt

START=99
STOP=01

OPTIONS="-vv"

SSD=start-stop-daemon

start() {
	$SSD -S -b -x hostile -- $OPTIONS
}

stop() {
	$SSD -K -p /var/run/hostile.pid
}
